# deploy/systemd/civmctl-runner-watchdog.timer
[Unit]
Description=Trigger civmctl-runner-watchdog.service a cada 2min
Documentation=https://github.com/advoq/civm/blob/main/docs/specs/civmctl/PRD.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=15s

[Install]
WantedBy=timers.target
