# deploy/systemd/civmctl-cleanup.timer
[Unit]
Description=Trigger civmctl-cleanup.service diariamente
Documentation=https://github.com/emersonbusson/ci-vm/blob/main/docs/specs/civmctl/PRD.md

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
