# deploy/systemd/civmctl-mem-watchdog.timer
[Unit]
Description=Trigger civmctl-mem-watchdog.service a cada 1min (monitoramento RAM near-real-time)
Documentation=https://github.com/advoq/civm/blob/main/internal/memwatchdog/

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
RandomizedDelaySec=5s

[Install]
WantedBy=timers.target
